# files/etc/uci-defaults/96-luci-zh-cn
#!/bin/sh
# Default LuCI language: Simplified Chinese

if uci -q get luci.main >/dev/null 2>&1; then
  uci set luci.main.lang='zh_cn'
  uci commit luci
fi

exit 0
